# files.backup/backup-OpenWrt-Prox-snap-2025-03-29/etc/uci-defaults/10_disable_services
/etc/init.d/adblock disable
/etc/init.d/gpio_switch disable
/etc/init.d/softethervpnclient disable

exit 0